# docker/s6/http-server/run
#!/bin/sh

su-exec app http-server /app
